# docker/sbdoc-with-mongo-alpine/supervisord.conf
[supervisord]
nodaemon=true

[program:mongodb]
command=/usr/bin/mongod

[program:sbdoc]
command=/usr/bin/node /root/SBDoc/SBDoc/bin/www
