# docker/ubuntu/supervisord.conf
[supervisord]
nodaemon=true

[program:DOClever]
command=node /root/DOClever/Server/bin/www
