# docker/supervisord.conf
[supervisord]
nodaemon=true

[program:mongodb]
command=/usr/bin/mongod

[program:DOClever]
command=node /root/DOClever/Server/bin/www
